# live-build/config/includes.chroot/etc/profile.d/ctx0an.sh
export CTX0AN_HOME=/opt/ctx0an
export PATH="/opt/ctx0an/venvs/ai-stack/bin:$PATH"