# live-build/config/includes.chroot/etc/profile.d/neuralcore.sh
export NEURALCORE_HOME=/opt/neuralcore
export PATH="/opt/neuralcore/venvs/ai-stack/bin:$PATH"
export OLLAMA_HOST=http://127.0.0.1:11434